# hallo/frame/supervisor/uwsgi.ini
[program:project]
command = uwsgi --ini uwsgi.ini
directory = /web/project
autostart = true
autorestart = true
redirect_stderr = true
stopasgroup = true
